# openwrt-packages/imagebuilder-overlay/files/etc/init.d/oxwrtd
#!/bin/sh /etc/rc.common
#
# procd init script for oxwrtd running as a side binary on a stock
# OpenWrt system. Starts in --control-only mode: only the sQUIC control
# plane runs; no mounts, no netlink, no firewall, no supervisor. Safe to
# leave enabled alongside procd/netifd/fw4 on an unmodified OpenWrt image.
#
# Operator workflow: `ssh root@<router> /etc/init.d/oxwrtd {start,stop,
# restart,enable,disable}`. `enable` plumbs it into procd's boot sequence
# so the control plane comes up on every reboot.
#
# When the firmware eventually takes over as PID 1 this script becomes
# redundant and will be dropped from the imagebuilder overlay in favor
# of the /sbin/init symlink path in openwrt-packages/oxwrtd/Makefile.

USE_PROCD=1

# Start late — after network + firewall are up. The control plane
# binds to 192.168.1.1:51820 which requires the LAN bridge to be
# addressed, so we can't start earlier than the network stack.
START=95
STOP=01

PROG=/usr/bin/oxwrtd

# Operator-selectable run mode. The single-line file /etc/oxwrt/mode
# determines how procd launches us:
#
#   control-only   (default if file missing) — supervisor empty, just
#                  the sQUIC control plane. Safe side-binary coexistence
#                  with stock OpenWrt — will not touch live network
#                  state even on reload.
#   services-only  supervisor populated from /etc/oxwrt/oxwrt.toml and spawns
#                  declared services (including debug-ssh on port 2222),
#                  but no network/firewall reconcile. Use this once
#                  you've validated the per-service rootfs layouts.
#   init           full PID-1-like takeover, including firewall install
#                  and netlink reconcile. DO NOT set this unless you're
#                  ready for the side-binary to fight procd/fw4 over
#                  /etc/config — this mode is really intended for when
#                  the binary IS actually PID 1 (via /sbin/init symlink).
#
# Anything else falls back to control-only with a log warning. Changing
# the mode means `ssh root@router /etc/init.d/oxwrtd restart`.
MODE_FILE=/etc/oxwrt/mode

get_mode() {
    if [ -r "$MODE_FILE" ]; then
        head -n1 "$MODE_FILE" | tr -d ' \t\n\r'
    else
        echo "control-only"
    fi
}

start_service() {
    # Refuse to start if the binary is missing — gives a clearer
    # failure than procd's generic "respawn too fast" loop.
    [ -x "$PROG" ] || {
        echo "oxwrtd: $PROG not found or not executable" >&2
        return 1
    }

    mode="$(get_mode)"
    case "$mode" in
        control-only|services-only|init)
            : ;;
        *)
            logger -t oxwrtd "unknown mode '$mode' in $MODE_FILE; falling back to control-only"
            mode="control-only"
            ;;
    esac
    logger -t oxwrtd "starting in --$mode mode"

    procd_open_instance oxwrtd
    procd_set_param command "$PROG" "--$mode"
    # Respawn on crash: threshold=3600s, timeout=5s, retry=5 times before
    # giving up. Matches procd defaults for long-running daemons.
    procd_set_param respawn 3600 5 5
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param pidfile /var/run/oxwrtd.pid
    # No file descriptors inherited — oxwrtd opens its own sockets.
    procd_set_param limits nofile=16384
    procd_close_instance
}

reload_service() {
    # oxwrtd's own `reload` RPC handles config reparsing in-band.
    # A procd-level reload = restart (simpler, guaranteed fresh state).
    stop
    start
}
