# openwrt-packages/imagebuilder-overlay/files/etc/uci-defaults/97-oxwrt-debug-ssh-rootfs
#!/bin/sh
# Build the debug-ssh service rootfs at /usr/lib/oxwrt/services/debug-ssh/rootfs
# on first boot.
#
# The "debug-ssh" service is an opt-in escape hatch for the PID-1
# takeover path: once oxwrtctl replaces procd, there is no stock
# dropbear and no shell, so a stuck control plane means a physical
# power-cycle + rollback flash. Running dropbear under our own
# supervisor on a second port keeps the SSH path available while
# staying under the full container hardening pipeline (caps drop,
# no_new_privs, pivot_root).
#
# The service is declared in oxwrt.toml but only spawns when oxwrtctl
# runs in --services-only or --init mode. It has no effect in
# --control-only mode (which is what /etc/init.d/oxwrtctl launches
# by default during bring-up).
#
# Port 2222 rather than 22 — stock dropbear still owns 22 while
# we're in side-binary mode. After PID-1 takeover this can move to 22.
#
# Runs before 98-oxwrt-diag-rootfs and 99-oxwrtctl (alphabetical) so
# the rootfs exists before the supervisor tries to spawn against it.
#
# Like the diag rootfs provisioner, this is only for the
# side-binary-on-stock-OpenWrt deployment. A production firmware
# either omits this rootfs entirely (appliance mode, no SSH) or
# ships a curated static-linked dropbear in its own OpenWrt package.

set -e
R=/usr/lib/oxwrt/services/debug-ssh/rootfs

# Two-stage idempotency: the binary + applet tree might be image-staged
# (imagebuilder-stage Makefile target copies them in at image-build
# time), but the dropbear host keys are per-device and regenerated on
# every clean flash. Check them separately so we always top up the key
# trio even if the binaries already exist.
need_bin=0
need_keys=0
[ -x "$R/sbin/dropbear" ] || need_bin=1
[ -s "$R/etc/dropbear/dropbear_ed25519_host_key" ] || need_keys=1

if [ "$need_bin" = 0 ] && [ "$need_keys" = 0 ]; then
    exit 0
fi

# Refuse to run if the host binary isn't present — a missing dropbear
# in the stock image means this script was enabled by mistake.
if [ ! -x /usr/sbin/dropbear ]; then
    logger -t oxwrt-debug-ssh "no /usr/sbin/dropbear on host; skipping"
    exit 0
fi

mkdir -p "$R/bin" "$R/sbin" "$R/lib" "$R/etc/dropbear" \
         "$R/dev/pts" "$R/proc" "$R/sys" "$R/tmp" "$R/var/log" "$R/root"

if [ "$need_bin" = 1 ]; then
    # Core binaries: busybox for the login shell (symlinked per-applet
    # below), dropbear for the listener itself.
    cp /bin/busybox "$R/bin/busybox"
    cp /usr/sbin/dropbear "$R/sbin/dropbear"

    # BusyBox applets a debug session actually uses. Adding more is
    # cheap (each is a symlink, ~10 bytes), but keeping the set visible
    # makes the rootfs auditable — "what can an attacker with SSH do?"
    # has a bounded answer.
    cd "$R/bin"
    for a in sh ash ls cat mount umount ps kill ip ping traceroute dmesg free \
             head tail grep sed awk echo chmod chown stat df du tr find wc \
             date touch mkdir rmdir rm cp mv ln readlink uname hostname \
             nslookup netstat which env printf id whoami pwd true false xargs \
             sleep clear reset login su tty who logger wget vi more less cut sort uniq; do
        ln -sf busybox "$a"
    done
    cd "$R/sbin"
    for a in ifconfig route arp; do
        ln -sf /bin/busybox "$a"
    done
    cd - >/dev/null

    # Runtime libs — dropbear and busybox both link against the musl
    # loader + libgcc_s on aarch64 OpenWrt.
    cp /lib/ld-musl-aarch64.so.1 "$R/lib/"
    cp /lib/libgcc_s.so.1 "$R/lib/"

    # Minimal /etc.
    cat > "$R/etc/passwd" <<EOF
root:x:0:0:root:/root:/bin/sh
EOF
    cat > "$R/etc/group" <<EOF
root:x:0:
EOF
    cat > "$R/etc/shells" <<EOF
/bin/sh
/bin/ash
EOF

    chmod 1777 "$R/tmp"
fi

# Host keys: generate directly via dropbearkey rather than copying
# from /etc/dropbear/.
#
# Earlier versions copied from stock dropbear's /etc/dropbear/. That's
# a race: uci-defaults runs in preinit, BEFORE /etc/init.d/dropbear
# starts — so on a fresh flash the host keys don't exist yet. Worse,
# after PID-1 takeover there IS no stock dropbear to generate them.
# Generating them here with dropbearkey (which ships alongside the
# dropbear binary) removes both failure modes.
#
# Keys are per-device: regenerated on every clean flash. If we want
# persistent keys across firmware updates, push them via the control
# plane post-flash — don't embed them in the image.
for kt in ed25519 rsa; do
    key="$R/etc/dropbear/dropbear_${kt}_host_key"
    if [ ! -s "$key" ]; then
        dropbearkey -t "$kt" -f "$key" >/dev/null 2>&1 || \
            logger -t oxwrt-debug-ssh "dropbearkey -t $kt failed"
    fi
done

# authorized_keys: ship the operator's key if stock dropbear has one
# staged (imagebuilder-stage copies ~/.ssh/id_ed25519.pub there). Safe
# to leave absent — the service will start but refuse logins, which
# is a better failure mode than not starting at all.
if [ -s /etc/dropbear/authorized_keys ]; then
    cp /etc/dropbear/authorized_keys "$R/etc/dropbear/" || true
fi

exit 0
