# openwrt-packages/imagebuilder-overlay/files/etc/uci-defaults/99-oxwrtd
#!/bin/sh
# Enable oxwrtd on first boot after flash.
#
# /etc/uci-defaults/* scripts run exactly once at first boot — procd
# deletes them afterward. This is the canonical OpenWrt pattern for
# "baked into the image but needs a boot-time action to take effect."
#
# We can't ship /etc/rc.d/S95oxwrtd directly because the imagebuilder
# overlay doesn't preserve symlinks reliably across filesystem hops.
# `/etc/init.d/oxwrtd enable` creates the symlink at first boot.

/etc/init.d/oxwrtd enable

exit 0
